# PreStage_Apps/pkgroot/Library/LaunchDaemons/com.github.smithjw.DEPNotify-prestarter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>com.github.smithjw.DEPNotify-prestarter</string>
        <key>RunAtLoad</key>
        <true/>
        <key>UserName</key>
        <string>root</string>
        <key>ProgramArguments</key>
        <array>
            <string>/bin/zsh</string>
            <string>/var/tmp/com.github.smithjw.DEPNotify-prestarter-installer.zsh</string>
        </array>
        <key>StandardErrorPath</key>
        <string>/var/tmp/com.github.smithjw.DEPNotify-prestarter-installer.zsh.err</string>
        <key>StandardOutPath</key>
        <string>/var/tmp/com.github.smithjw.DEPNotify-prestarter-installer.zsh.out</string>
    </dict>
</plist>